# supervisor.conf
[program:tg_resend]
command=/www/tg_resend_venv/bin/python main.py
directory=/www/tg_resend/
user=www-data
autorestart=true
